# deploy/hetzner/inference-bazaar-arti.service
[Unit]
Description=InferenceBazaar venue Tor onion service (Arti)
After=network-online.target
Wants=network-online.target

[Service]
# Runs Arti as the venue's onion service (see arti.toml). Install Arti with
# `cargo install arti` or the distro package; adjust ExecStart to its path.
ExecStartPre=/usr/bin/install -d -o inference-bazaar -g inference-bazaar /opt/inference-bazaar/arti/state /opt/inference-bazaar/arti/cache
ExecStart=/usr/local/bin/arti -c /opt/inference-bazaar/arti.toml proxy
# Publish the generated .onion to the stable path the operator reads each time
# Arti (re)writes it, so a fresh keypair is surfaced without a venue restart.
ExecStartPost=/opt/inference-bazaar/publish-onion.sh
Restart=always
RestartSec=5
User=inference-bazaar
Group=inference-bazaar

[Install]
WantedBy=multi-user.target
