# deploy/hetzner/surplus3-blueprint-runtime.service
[Unit]
Description=Surplus venue 2 (Base Sepolia blueprint 17 service 3, operator 0x2420)
After=network.target surplus-mm-sidecar.service
Wants=surplus-mm-sidecar.service

[Service]
Environment=RUST_LOG=info,surplus_operator=debug
Environment=HTTP_RPC_URL=https://base-sepolia-rpc.publicnode.com/
Environment=WS_RPC_URL=wss://base-sepolia.drpc.org/
Environment=KEYSTORE_URI=/opt/surplus3/keystore
Environment=DATA_DIR=/opt/surplus3/data
Environment=BLUEPRINT_ID=17
Environment=SERVICE_ID=3
Environment=TANGLE_CONTRACT=0x8299d60f373f3a4a8c4878e335cb9d840e6e3730
Environment=STAKING_CONTRACT=0x91b1186f4f31d6e02e481c0af29c7244a3fe417d
Environment=STATUS_REGISTRY_CONTRACT=0x2a7ceb96a9b18721b5bbb0022b4d358b3c50bcb2
Environment=PROTOCOL=tangle
Environment=TEST_MODE=true
Environment=SURPLUS_SIDECAR_URL=http://127.0.0.1:9310
Environment=SURPLUS_ROUTER_URL=https://router.tangle.tools
# Bonded issuer: set a real backend (router-proxy mode is refused at boot).
Environment=SURPLUS_INFERENCE_URL=<issuer-openai-compatible-backend>
Environment=SURPLUS_OPERATOR_ADDR=127.0.0.1:9500
Environment=SURPLUS_MM_SIZE=1500000
Environment=SURPLUS_MM_LEVELS=3
Environment=SURPLUS_CHAIN_ID=84532
Environment=SURPLUS_SETTLEMENT_ADDR=0x3fa622488fD970ECdE23b8384a98de6fFa5A1763
# Attester / co-sign identity (address must be in the book's on-chain attesters).
Environment=SURPLUS_OPERATOR_KEY=<issuer-attester-key>
# Separate gas/submitter key (key separation; falls back to operator key if unset).
Environment=SURPLUS_SUBMITTER_KEY=<gas-submitter-key>
Environment=SURPLUS_RPC_URL=https://sepolia.base.org
# Optional Tor onion endpoint (run Arti + expose /redeem as an onion service);
# published via /health so privacy-mode clients dial it. Leave unset if no onion.
# Environment=SURPLUS_ONION_URL=http://<your-onion>.onion
# Environment=PRIVACY_MODE=tor   # route the operator's OWN outbound calls via Arti
ExecStart=/opt/surplus/surplus-operator run
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
